# systemd/system/yazelix_host_policy.path
[Unit]
Description=Watch profile-owned FlexNetOS host cache policy inputs

[Path]
PathChanged=/etc/nix/nix.conf
PathChanged=/etc/nix/nix.custom.conf
PathChanged=/etc/determinate/config.json
PathChanged=/etc/shells
Unit=yazelix_host_policy.service

[Install]
WantedBy=multi-user.target
