# systemd/system/nix-daemon.service.d/10-yazelix-host-policy.conf
[Unit]
Requires=yazelix_host_policy.service
After=yazelix_host_policy.service

[Service]
ExecStart=
ExecStart=@/home/flexnetos/.nix-profile/bin/nix-daemon nix-daemon --daemon
