# systemd/user/flexnetos_runner@.service
[Unit]
Description=FlexNetOS ephemeral profile-owned runner %i
After=network-online.target
Wants=network-online.target
ConditionPathIsExecutable=/home/flexnetos/.nix-profile/toolbin/nu
ConditionPathIsExecutable=/home/flexnetos/.nix-profile/bin/kache
ConditionPathIsExecutable=/home/flexnetos/.nix-profile/bin/fxrun-actions

[Service]
Type=exec
Environment=SHELL=/home/flexnetos/.nix-profile/toolbin/nu
Environment=HOME=/run/user/1001/yazelix/runners/%i/home
Environment=XDG_DATA_HOME=/home/flexnetos/meta/var/lib
Environment=XDG_STATE_HOME=/home/flexnetos/meta/var/lib
Environment=XDG_CACHE_HOME=/run/user/1001/yazelix/runners/%i/xdg_cache
Environment=YAZELIX_STATE_DIR=/run/user/1001/yazelix/runners/%i/yazelix
Environment=CARGO_HOME=/run/user/1001/yazelix/runners/%i/cargo_home
Environment=CARGO_TARGET_DIR=/run/user/1001/yazelix/runners/%i/cargo_target
Environment=KACHE_BIN=/home/flexnetos/.nix-profile/bin/kache
Environment=KACHE_CACHE_DIR=/home/flexnetos/.cache/kache/runners/%i
Environment=RUSTC_WRAPPER=/home/flexnetos/.nix-profile/bin/kache-rustc-wrapper
Environment=CARGO_BUILD_RUSTC_WRAPPER=/home/flexnetos/.nix-profile/bin/kache-rustc-wrapper
Environment=GH_CONFIG_DIR=/home/flexnetos/.config/gh
Environment=CODEX_HOME=/run/user/1001/yazelix/profile-runtime/codex
Environment=CLAUDE_CONFIG_DIR=/home/flexnetos/meta/var/lib/claude
ExecStartPre=/home/flexnetos/.nix-profile/bin/flexnetos_runner_policy runtime %i
ExecStart=/home/flexnetos/.nix-profile/bin/flexnetos_runner_service %i
StandardOutput=null
StandardError=null
# Per-job ephemeral (each job gets a fresh JIT runner), but the service is
# persistent: after a job the runner exits and systemd restarts it, re-running
# ExecStartPre to re-materialize the runner tree and register a new JIT runner.
# This is the GitHub-recommended ephemeral-plus-always-restart pattern and keeps
# the fork-PR isolation posture (one job per runner registration).
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
